# keyboard-colors.service
[Unit]
Description=System76 Keyboard Colors
[Service]
Type=Simple
ExecStart=/usr/bin/php /opt/keyboard-colors/keyboard-color.php --frequency 2 --mode rainbow # update every 2 seconds
[Install]
WantedBy=multi-user.target