# keyboard-colors.service
[Unit]
Description=System76 Keyboard Colors
[Service]
Type=simple
ExecStart=/usr/local/bin/keyboard-color -c /etc/keyboard-color.json
PIDFile=keyboard-colors.pid
[Install]
WantedBy=multi-user.target
